# gdrive-backup.timer
[Unit]
Description=Run Google Drive backup daily

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target